# vfd-controller.service
[Unit]
Description=TX3 Mini VFD Controller
Documentation=https://github.com/marufmoinuddin/tx3-mini-vfd
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/vfd_controller
Restart=always
RestartSec=5
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/sys/devices/platform/spi /sys/class/leds
ProtectHome=true
ProtectKernelTunables=false
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
MemoryLimit=50M
CPUQuota=10%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vfd-controller

[Install]
WantedBy=multi-user.target
